# deploy/ipam.service
# ==============================================================================
# IPAM SIW v2 — Service systemd pour Node.js
# Fichier : /etc/systemd/system/ipam.service
#
# Installation :
#   cp ipam.service /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable ipam
#   systemctl start ipam
# ==============================================================================

[Unit]
Description=IPAM SIW v2 — Serveur Node.js (Redis)
After=network.target redis.service
Wants=redis.service

[Service]
WorkingDirectory=/var/www/ipam

ExecStart=/usr/bin/node /var/www/ipam/server/index.mjs

Restart=always
RestartSec=10

Environment=NODE_ENV=production
Environment=PORT=3000
Environment=REDIS_URL=redis://127.0.0.1:6379

# ── Sécurité ──────────────────────────────────────────────────────────────────
# Le service tourne en root pour pouvoir appeler systemctl/journalctl/reboot
# PrivateTmp retiré : isolait /tmp, bloquait l'explorateur de fichiers
# RestrictAddressFamilies retiré : bloquait AF_NETLINK (ip addr) et AF_PACKET (tcpdump)
ProtectSystem=full
ReadWritePaths=/var/www/ipam/data /var/lib/redis
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
MemoryMax=512M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ipam

[Install]
WantedBy=multi-user.target
